# roles/fancontrol/files/fancontrol-restart.service
[Unit]
After=suspend.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart fancontrol

[Install]
WantedBy=suspend.target
